# rpi/sourdough.service
[Unit]
Description=Sourdough Data Collector
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/sourdough
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Fill these in on the Pi (e.g. via /etc/systemd/system/sourdough.service.d/override.conf)
Environment=SUPABASE_URL=
Environment=SUPABASE_SERVICE_KEY=
Environment=MEASUREMENT_INTERVAL=300
Environment=SYNC_INTERVAL=900
Environment=DB_PATH=/home/pi/sourdough/data/sourdough.db

[Install]
WantedBy=multi-user.target
